# supervisor/wsl2/conf.d/wsl22/youtube-react.conf
[program:youtube_react]
environment=PATH=%(ENV_PATH)s:/usr/local/node-v18.12.1-linux-x64/bin
directory=/home/utylee/.virtualenvs/react/src/youtube-uploader/
# env PATH에 node 바이너리 위치를 추가해주고 npm을 실행하면 되는 것 같습니다
# command=env PATH=%(ENV_PATH)s:/usr/local/node-v18.12.1-linux-x64/bin /usr/local/node-v18.12.1-linux-x64/bin/npm run start
#command=npm run start
command=/usr/local/node-v18.12.1-linux-x64/bin/npm run start
stopsignal=KILL
stopasgroup=true
# 리빌드 안하고 시작할 경우 계속 시작하며 프로세스가 멈추는 경우가 많았습니다
# 안타깝게도 기대한대로 동작이 안되는 것 같습니다
autorestart=unexpected
startretries=3
#killasgroup=true
#user=nobody
user=root
#user=utylee
