# supervisor/wsl1/conf.d/capture_watcher.conf
[program:capture_watcher]
command = /home/utylee/.virtualenvs/misc/bin/python /home/utylee/.virtualenvs/misc/src/capture_watcher.py
autostart = true
autorestart = true
