# supervisor/wsl1/conf.d/upload_backend.conf
[program:api_youtube_upload]
command = /home/utylee/.virtualenvs/api/bin/python /home/utylee/.virtualenvs/api/src/api_youtube.py
autostart = true
autorestart = true
