# supervisor/wsl2/conf.d/wsl21/cookie_changer.conf
[program:cookie_changer]
# environment=env PATH=%(ENV_PATH)s:/usr/local/node-v18.12.1-linux-x64/bin
environment=PATH=%(ENV_PATH)s:/mnt/c/Windows/System32/WindowsPowerShell/v1.0/
# directory=/home/utylee/.virtualenvs/react/src/youtube-uploader/
# env PATH에 node 바이너리 위치를 추가해주고 npm을 실행하면 되는 것 같습니다
#command=env PATH=%(ENV_PATH)s:/usr/local/node-v18.12.1-linux-x64/bin /usr/local/node-v18.12.1-linux-x64/bin/npm run start
#
# command=/home/utylee/.virtualenvs/misc/bin/python /home/utylee/.virtualenvs/misc/src/cookie_changer.py --path /tmp/cookie_changer.sock
command=/home/utylee/.virtualenvs/misc/bin/python /home/utylee/.virtualenvs/misc/src/cookie_changer.py --port 9978
user=utylee

